# sphinx/docker/source/_static/code/runtime/ng/nginx.ini
[program:nginx]
command=nginx -g "daemon off;"
priority=999