# sphinx/docker/source/_static/code/arm/rpi-jupyter/jupyter.conf
[program:jupyter]
command=jupyter %(ENV_JUPYTER_TYPE)s --ip='0.0.0.0' --NotebookApp.token='' --NotebookApp.password='' --notebook-dir=/ipynb --no-browser --allow-root
priority=999