# docker/aws-developer/ubuntu/etc/supervisor/conf.d/aws-cli.conf
[program:awscli]
command=/usr/local/sbin/setup-awscli.sh 
startsecs=0
exitcodes=0
priority=999