# docker/spark-intro/ubuntu/etc/supervisor/conf.d/streamalpha.conf
[program:streamalpha]
command=/usr/local/conda/bin/python /usr/local/bin/stream-alpha.py
exitcodes=0
startsecs=0
priority=999