# docker/spark-intro/ubuntu/etc/supervisor/conf.d/streamnum.conf
[program:streamnum]
command=/usr/local/conda/bin/python /usr/local/bin/stream-num.py
exitcodes=0
startsecs=0
priority=999