# spark-intro/container/ubuntu/etc/supervisor/conf.d/all.conf
[program:all]
command=/bin/bash /usr/local/bin/start-all.sh
exitcodes=0
startsecs=0
priority=999